# roles/tomcat/templates/tomcat-supervisor.ini

[program:tomcat-{{app}}]
command = {{ deploy_path }}/tomcat-{{app}}/bin/catalina.sh run
autorestart=true
redirect_stderr = true
stdout_logfile = {{ deploy_path }}/tomcat-{{app}}/logs/catalina.out
logfile_maxbytes = 10MB
user = {{ tomcat_user }}
