# roles/tomcat/templates/tomcat-supervisor.ini

[program:tomcat-{{app}}]
command = {{ deploy_path }}/%(program_name)s/bin/catalina.sh run
autorestart=true
redirect_stderr = true
stdout_logfile = {{ deploy_path }}/%(program_name)s/logs/catalina.out
logfile_maxbytes = 10MB
user = {{ tomcat_user }}
